# native/launchd/io.hexmux.supervisor.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.hexmux.supervisor</string>
  <key>ProgramArguments</key>
  <array>
    <string>@HEXMUX_ACTIVATE@</string>
    <string>--socket-name</string>
    <string>Listeners</string>
    <string>--</string>
    <string>@PYTHON@</string>
    <string>-m</string>
    <string>hexmux.supervisor</string>
    <string>--idle-timeout</string>
    <string>30</string>
  </array>
  <key>Sockets</key>
  <dict>
    <key>Listeners</key>
    <dict>
      <key>SockPathName</key>
      <string>@SOCKET_PATH@</string>
      <key>SockPathMode</key>
      <integer>384</integer>
      <key>SockType</key>
      <string>stream</string>
    </dict>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
